# etc/systemd/system/spelunkicons.service
[Unit]
Description=Spelunkicons

[Service]
TimeoutStartSec=0
Restart=always
ExecStart=/opt/spelunkicons

[Install]
WantedBy=default.target
